# smartbot.service
[Unit]
Description=SmartBot AI VK Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/smartbot
ExecStart=/usr/bin/python3 /root/smartbot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Переменные окружения
Environment=PYTHONPATH=/root/smartbot
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
